# linux/systemd/.config/systemd/user/spotifyd.service
[Unit]
Description=spotifyd

[Service]
ExecStart=/bin/spotifyd --no-daemon --use-keyring
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=default.target
